# runner/ekos-runner.service
[Unit]
Description=Observatory EKOS Runner Agent
Documentation=https://github.com/UH-Bayfordbury/Observatory-Automation
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Run as a dedicated user with access to the KStars D-Bus session if real hardware is used.
User=astro
ExecStart=/usr/local/bin/ekos-runner -config /etc/ekos-runner/runner.yaml
Restart=on-failure
RestartSec=5
# State directory for generated EKOS bundles (-> /var/lib/ekos-runner).
StateDirectory=ekos-runner
# Harden a little.
NoNewPrivileges=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
